# launchd/bots/com.antigravity.seonyuduo-exercise-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.seonyuduo-exercise-bot</string>

    <key>UserName</key>
    <string>__MACMINI_USER__</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/__MACMINI_USER__/Antigravity_Market_Dashboard/launchd/bots/run_bot.sh</string>
        <string>seonyuduo-exercise-bot</string>
        <string>/Users/__MACMINI_USER__/Antigravity_Market_Dashboard/venv/bin/python3</string>
        <string>/Users/__MACMINI_USER__/Antigravity_Market_Dashboard/execution/seonyuduo_exercise_bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/__MACMINI_USER__/Antigravity_Market_Dashboard/execution</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/__MACMINI_USER__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>PYTHONIOENCODING</key>
        <string>utf-8</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/__MACMINI_USER__/Antigravity_Market_Dashboard/logs/launchd/seonyuduo-exercise-bot.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/__MACMINI_USER__/Antigravity_Market_Dashboard/logs/launchd/seonyuduo-exercise-bot.err</string>
</dict>
</plist>
