# datalake/launchd/com.antigravity.datalake-webui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.datalake-webui</string>
    <key>UserName</key>
    <string>__MACMINI_USER__</string>
    <key>WorkingDirectory</key>
    <string>__REPO__</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/datalake/webui/run_webui.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/__MACMINI_USER__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>
    <!-- 상시 데몬: 127.0.0.1:8787 — Caddy /wiki/* 프록시의 백엔드 (AoE Wiki 탭, 2026-07-12) -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>
    <key>StandardOutPath</key>
    <string>__REPO__/logs/launchd/datalake-webui.out</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/launchd/datalake-webui.err</string>
</dict>
</plist>
