# datalake/launchd/com.antigravity.datalake-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.datalake-backup</string>
    <key>UserName</key>
    <string>__MACMINI_USER__</string>
    <key>WorkingDirectory</key>
    <string>__REPO__</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/datalake/launchd/run_datalake_job.sh</string>
        <string>datalake-backup</string>
    </array>
    <!-- 데몬 환경 표준 (CONTRACT v3): launchd 는 로그인 환경이 없음 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/__MACMINI_USER__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>
    <!-- KST 일요일 10:00 — datalake 백업 push (시스템 TZ=Asia/Seoul 전제) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__REPO__/logs/launchd/datalake-backup.out</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/launchd/datalake-backup.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
