# launchd/system/com.antigravity.git-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.git-pull</string>

    <!-- Sync the repo every 5 minutes (VM */5 git-pull cron replacement).
         RunAtLoad syncs once shortly after boot too; a boot-time failure (no
         network yet) is harmless — it just increments the streak, which the
         next successful pull resets. -->
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/launchd/system/git_pull.sh</string>
    </array>

    <!-- Run as the migration user so the checkout + pull use that account's git
         config / credentials, and log/state files are user-owned. -->
    <key>UserName</key>
    <string>__MACMINI_USER__</string>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- launchd has no login environment; CONTRACT 데몬 환경 표준. HOME is
         required for git to find ~/.gitconfig and credentials. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/__MACMINI_USER__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>TZ</key>
        <string>Asia/Seoul</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__REPO__/logs/launchd/git-pull.out</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/launchd/git-pull.err</string>

    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
