# launchd/web/com.antigravity.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key><string>com.antigravity.web</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/caddy</string>
		<string>run</string>
		<string>--config</string>
		<string>/Users/sisyphe/Antigravity_Market_Dashboard/launchd/web/Caddyfile</string>
		<string>--adapter</string>
		<string>caddyfile</string>
	</array>
	<key>UserName</key><string>sisyphe</string>
	<key>RunAtLoad</key><true/>
	<key>KeepAlive</key><true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key><string>/Users/sisyphe</string>
		<key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin</string>
	</dict>
	<key>StandardOutPath</key><string>/Users/sisyphe/Antigravity_Market_Dashboard/logs/launchd/web.out</string>
	<key>StandardErrorPath</key><string>/Users/sisyphe/Antigravity_Market_Dashboard/logs/launchd/web.err</string>
</dict>
</plist>
