# launchd/web/com.antigravity.sisyphe-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key><string>com.antigravity.sisyphe-pull</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/git</string>
		<string>-C</string>
		<string>/Users/sisyphe/srv/sisyphe_repo</string>
		<string>pull</string>
		<string>--ff-only</string>
		<string>--quiet</string>
	</array>
	<key>UserName</key><string>sisyphe</string>
	<key>StartInterval</key><integer>3600</integer>
	<key>RunAtLoad</key><true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key><string>/Users/sisyphe</string>
		<key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin</string>
	</dict>
	<key>StandardOutPath</key><string>/Users/sisyphe/Antigravity_Market_Dashboard/logs/launchd/sisyphe-pull.out</string>
	<key>StandardErrorPath</key><string>/Users/sisyphe/Antigravity_Market_Dashboard/logs/launchd/sisyphe-pull.err</string>
</dict>
</plist>
