# scripts/earnings-bot-notify.service
[Unit]
Description=Earnings Bot Failure Notification

[Service]
Type=oneshot
ExecStart=/home/ubuntu/Antigravity_Market_Dashboard/scripts/notify_sisyphe_failure.sh earnings-bot
