# scripts/earnings-bot.timer
[Unit]
Description=Antigravity Earnings Bot timer (매일 08:00 KST 1회)

[Timer]
# 매일 한국 시간 08:00 (UTC 23:00 전날) — 한 번 호출에 캘린더 sync + 폴링 + 분석 + 발행 + 다이제스트 모두 진행
OnCalendar=*-*-* 08:00:00 Asia/Seoul
AccuracySec=1min
Unit=earnings-bot.service
Persistent=true

[Install]
WantedBy=timers.target
