# scripts/hotel-adr.service
[Unit]
Description=Antigravity Hotel ADR Fetch (booking.com 호텔 ADR lead +7/+14/+30)
After=network-online.target
Wants=network-online.target
OnFailure=sisyphe-bot-notify@hotel-adr.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/Antigravity_Market_Dashboard
ExecStart=/bin/bash /home/ubuntu/Antigravity_Market_Dashboard/scripts/run_hotel_adr.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=20min

[Install]
WantedBy=multi-user.target
