# scripts/hotel-adr.service
[Unit]
Description=Antigravity Hotel ADR Fetch (booking.com 호텔 ADR lead +7/+14/+30)
After=network-online.target
Wants=network-online.target
OnFailure=sisyphe-bot-notify@hotel-adr.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/Antigravity_Market_Dashboard
ExecStart=/bin/bash /home/ubuntu/Antigravity_Market_Dashboard/scripts/run_hotel_adr.sh
StandardOutput=journal
StandardError=journal
# booking.com 30 페이지 × ~37s/page(eager) ≈ 18분 → 35분 캡(넉넉히). 일 1회 백그라운드.
TimeoutStartSec=35min

[Install]
WantedBy=multi-user.target
