# scripts/landing-highlights.timer
[Unit]
Description=Antigravity Landing Highlights timer (매일 18:45 KST, Featured 2차 마감 후 — 2026-07-07 18:35→18:45 이격(Featured 타임아웃 경합 해소, audits/schedule_conflicts.md R1))

[Timer]
OnCalendar=*-*-* 18:45:00 Asia/Seoul
RandomizedDelaySec=60
AccuracySec=30s
Unit=landing-highlights.service
Persistent=true

[Install]
WantedBy=timers.target
