# scripts/sisyphe-bot-notify@.service
[Unit]
Description=Sisyphe-family failure notification (for %i)

[Service]
Type=oneshot
ExecStart=/home/ubuntu/Antigravity_Market_Dashboard/scripts/notify_sisyphe_failure.sh %i
