# launchd/timers/com.antigravity.send-advisory-emails.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.send-advisory-emails</string>
    <key>UserName</key>
    <string>__MACMINI_USER__</string>
    <key>WorkingDirectory</key>
    <string>__REPO__</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/launchd/timers/run_timer_job.sh</string>
        <string>send-advisory-emails</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/__MACMINI_USER__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>
    <!-- 60초 폴러: StartCalendarInterval 대신 StartInterval. 잡은 요청 없으면 즉시 exit0(락으로 중복 방어) -->
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>__REPO__/logs/launchd/send-advisory-emails.out</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/launchd/send-advisory-emails.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
